# meta-tyan/meta-common/recipes-userspace/host-power-ctrl/files/service-files/tyan-host-reset.service
[Unit]
Description=Cold Reset Host
Requires=tyan-host-ready.target
After=tyan-host-ready.target

[Service]
Type=oneshot
ExecStart=/usr/bin/host_reset.sh
